# setup/supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:sshd]
command=/usr/sbin/sshd -D

[program:startxvfb]
command=/bin/sh /root/startXvfb.sh
